# config/autostart/elephant.desktop
[Desktop Entry]
Exec=elephant
Icon=
Name=elephant
Path=
Terminal=true
Type=Application
